# ai-helper.service
[Unit]
Description=AI Helper Application
After=network.target
Wants=network.target

[Service]
Type=simple
User=yapayzeka
Group=yapayzeka
WorkingDirectory=/opt/ai-helper
Environment=PATH=/opt/ai-helper/venv/bin
ExecStart=/opt/ai-helper/start_service.sh
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
